# deployment/qa/3030.iguzman-db.qa.supervisor.conf
[program:iguzman-db-qa]
directory=/var/www/apps/qa/iguzman-db/
environment=env=qa,db_name=iguzman_qa,db_user=iguzman_qa,db_password=iguzman_qa,worker=0
command=/opt/python-env/iguzman-qa/bin/gunicorn --bind 0.0.0.0:3030 --workers 2 --reload iguzman_db.wsgi:application
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/iguzman-db.qa.error.log
stdout_logfile=/var/log/supervisor/iguzman-db.qa.out.log
