# deployment/master/3034.iguzman-db.master.supervisor.conf
[program:iguzman-db-master-3034]
directory=/var/www/apps/master/iguzman-db/
environment=env=master,db_name=iguzman_master,db_user=iguzman_master,db_password=iguzman_master,worker=3
command=/opt/python-env/iguzman-master/bin/gunicorn --bind 0.0.0.0:3034 --workers 2 --reload iguzman_db.wsgi:application
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/iguzman-db.master.3034.error.log
stdout_logfile=/var/log/supervisor/iguzman-db.master.3034.out.log
