# deployment/staging/3031.iguzman-db.staging.supervisor.conf
[program:iguzman-db-staging]
directory=/var/www/apps/staging/iguzman-db/
environment=env=staging,db_name=iguzman_staging,db_user=iguzman_staging,db_password=iguzman_staging,worker=0
command=/opt/python-env/iguzman-staging/bin/gunicorn --bind 0.0.0.0:3031 --workers 2 --reload iguzman_db.wsgi:application
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/iguzman-db.staging.error.log
stdout_logfile=/var/log/supervisor/iguzman-db.staging.out.log
